# local/services/polkit-gnome-authentication-agent-1/run
#!/usr/bin/env dash
#
# João Farias © 2024 BeyondMagic <beyondmagic@mail.ru>

exec /usr/lib/polkit-gnome/polkit-gnome-authentication-agent-1
