# local/services/pueued/finish
#!/usr/bin/env nu
#
# João Farias © 2024 BeyondMagic <beyondmagic@mail.ru>

use service_manager.nu

def main [number: int, number2: int] -> int {
	service_manager finish $number 'pueued' "task manager"
}
