# local/services/ssh-agent/run
#!/usr/bin/env dash
#
# João Farias © 2024 BeyondMagic <beyondmagic@mail.ru>

exec ssh-agent -D -a /tmp/ssh-agent.socket
