# local/services/cliphist/run
#!/usr/bin/env dash
#
# João Farias © 2024 BeyondMagic <beyondmagic@mail.ru>

exec wl-paste --watch cliphist store
