# run/runit/service/runsvdir-dream/run
#!/usr/bin/env nu
#
# João Farias © 2024 BeyondMagic <beyondmagic@mail.ru>

const user = 'dream'
const home = '/home/dream/'
$env.USER = $user
$env.HOME = $home

source-env ($home | path join .config/nushell/scripts/environment.nu)

let groups = (id -Gn $user | tr ' ' ':')
let svdir = $"($home)/.local/services/"

exec chpst -u $"($user):($groups)" runsvdir $svdir
